# ansible/roles/goapp/templates/goapp.service
[Unit]
Description=Go App HTTP Example
After=network.target

[Service]
ExecStart={{goAppDir}}/app

[Install]
WantedBy=multi-user.target
